# supervisord.conf
[supervisord]
nodaemon=true

[program:XXX]
command=java -jar XXX


[inet_http_server]
port = *:9001
username = supervisord
password = 123